# .config/autostart/imwheel.sh.desktop
[Desktop Entry]
Exec=/home/tbs/.config/autostart/imwheel.sh
Icon=dialog-scripts
Name=imwheel.sh
Path=
Type=Application
X-KDE-AutostartScript=true
